# sites/studio/deploy/cloudflared.service
[Unit]
Description=Cloudflare Tunnel for Inertia Studio
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=agent
ExecStart=/usr/local/bin/cloudflared tunnel --config /home/agent/.cloudflared/config.yml run
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudflared

[Install]
WantedBy=multi-user.target
